# pkg/system/templates/etc/systemd/system/kubelet.service.d/20-max-pods.conf
[Service]
Environment='KUBELET_MAX_PODS=--max-pods={{.Node.MaxPods}}'
