# pkg/system/templates/etc/systemd/system/kubelet.service.d/30-kube-reserved.conf
[Service]
Environment='KUBELET_KUBE_RESERVED=--kube-reserved=cpu={{.Node.ReservedCPU}},memory={{.Node.ReservedMemory}}'
