# pkg/system/templates/etc/systemd/system/kubelet.service.d/30-taints.conf
[Service]
{{- if .Node.Taints }}
Environment='KUBELET_NODE_TAINTS=--register-with-taints="{{ range $index, $taint := .Node.Taints }}{{ if $index }},{{ end }}{{ $taint }}{{ end }}"'
{{ end -}}
