# pkg/system/templates/etc/systemd/system/kubelet.service.d/40-labels.conf
[Service]
{{- if .Node.Labels }}
Environment='KUBELET_NODE_LABELS=--node-labels="{{ range $index, $label := .Node.Labels }}{{ if $index }},{{ end }}{{ $label }}{{ end }}"'
{{ end -}}
